# deploy/litestream.service
# /etc/systemd/system/litestream.service
# Continuous SQLite replication to R2. Runs alongside wanderwallet.service.
[Unit]
Description=Litestream (continuous SQLite replication to R2)
After=network.target wanderwallet.service

[Service]
Type=simple
User=wanderwallet
EnvironmentFile=/opt/wanderwallet/deploy/backup.env
ExecStart=/usr/bin/litestream replicate -config /opt/wanderwallet/deploy/litestream.yml
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
