# deploy/wanderwallet.service
# /etc/systemd/system/wanderwallet.service
# Runs the production Next.js server (UI + API + webhook + in-process worker)
# as a non-root user. Deploy steps: build first (`npm ci && npm run build`),
# then `systemctl enable --now wanderwallet`.
[Unit]
Description=Wanderwallet (Next.js)
After=network.target

[Service]
Type=simple
User=wanderwallet
WorkingDirectory=/opt/wanderwallet
EnvironmentFile=/opt/wanderwallet/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3
# graceful shutdown: let in-flight jobs finish
KillSignal=SIGTERM
TimeoutStopSec=30
# memory backstop; tune to the VPS
MemoryMax=512M
# logs go to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
