# deploy/wanderwallet-backup.service
# /etc/systemd/system/wanderwallet-backup.service
# One-shot uploads backup, triggered by wanderwallet-backup.timer.
[Unit]
Description=Wanderwallet uploads backup (restic to R2)

[Service]
Type=oneshot
User=wanderwallet
WorkingDirectory=/opt/wanderwallet
ExecStart=/opt/wanderwallet/deploy/backup-uploads.sh
